# deploy/quadlet/thingino-broker.container
# Quadlet unit (rootful): /etc/containers/systemd/thingino-broker.container
# @DIR@ is substituted with the repo path by deploy.sh.
[Unit]
Description=Thingino web-builder broker
After=network-online.target
Wants=network-online.target

[Container]
# Published by CI to ghcr.io; @IMAGE@ is substituted (image:tag) by deploy.sh.
Image=@IMAGE@
# Let `podman auto-update` swap in a newer image from the registry (UI-triggered).
AutoUpdate=registry
EnvironmentFile=@DIR@/.env
Environment=IP_HEADER=X-Real-IP
Volume=thingino-broker-data:/data
# GitHub App auth (instead of GITHUB_TOKEN): drop app-key.pem in the repo dir,
# uncomment the mount below, and set GITHUB_APP_KEY_PATH=/app/app-key.pem in .env.
#Volume=@DIR@/app-key.pem:/app/app-key.pem:ro,Z
# Loopback only — reached by the host-networked Caddy.
PublishPort=127.0.0.1:8080:8080

[Service]
Restart=always
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
